# crontab
# START CRON JOB
0 23 * * * /usr/local/bin/python3 /app/ImmichAutoAlbum.py > /proc/1/fd/1 2>/proc/1/fd/2 
# END CRON JOB